# com.kian.transcribe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        Template plist for COO local automation.
        Copy this file and replace placeholders before loading:
        - __COO_LABEL__
        - __COO_RUNNER_SCRIPT__
        - __VOICE_MEMOS_DIR__
        - __WORKING_DIR__

        This watcher triggers COO processing whenever VOICE_MEMOS_DIR changes.
    -->
    <key>Label</key>
    <string>__COO_LABEL__</string>
    <key>ProgramArguments</key>
    <array>
        <string>__COO_RUNNER_SCRIPT__</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>__VOICE_MEMOS_DIR__</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__WORKING_DIR__</string>
</dict>
</plist>
